# template_master.ini
[program:{{NAME}}]
directory = {{WEED_DIR}}
command = {{WEED_BIN}} master -port={{M_PORT}} -ip={{NODE}} -mdir={{DATA}} -peers={{PEERS}} -defaultReplication={{REPLICATION}}
autostart = true
startsecs = 5
autorestart = true
user = root
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 10
stdout_logfile = /etc/supervisord.d/log/{{NAME}}.log